# deploy/monster-solana.service
[Unit]
Description=Monster Solana Predator
After=network.target

[Service]
Type=simple
User=monster
WorkingDirectory=/opt/monster-scanner
Environment=PATH=/opt/monster-scanner/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/monster-scanner/.venv/bin/python solana_predator.py
Restart=always
RestartSec=30
StandardOutput=append:/opt/monster-scanner/logs/solana_predator.log
StandardError=append:/opt/monster-scanner/logs/solana_predator.log

[Install]
WantedBy=multi-user.target
