# deploy/monster-ens-sniper.service
[Unit]
Description=Monster ENS Sniper
After=network.target

[Service]
Type=simple
User=monster
WorkingDirectory=/opt/monster-scanner
Environment=PATH=/opt/monster-scanner/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/monster-scanner/.venv/bin/python ens_sniper.py
Restart=always
RestartSec=30
StandardOutput=append:/opt/monster-scanner/logs/ens_sniper.log
StandardError=append:/opt/monster-scanner/logs/ens_sniper.log

[Install]
WantedBy=multi-user.target
